# roles/vnc/templates/novnc.j2.conf
description     "start noVNC websocket proxy"

start on started

exec /usr/local/bin/novnc.sh --listen {{ vnc_listen }} --vnc {{ vnc_host }}:{{ vnc_port }}
