# roles/vnc/templates/vncserver.j2.conf
description     "start vncserver as developer"

start on started

setuid {{ vnc_user }}

post-start script
  exec /usr/bin/vncserver :1
end script
