# roles/vnc/templates/vncserver.j2.conf
description     "start vncserver as {{ vnc_user }}"

start on started

post-start script
  exec su {{ vnc_user }} -c '/usr/bin/vncserver :1'
end script
