# roles/vnc/templates/vncserver.j2.conf
description     "start vncserver as {{ vnc_user }}"

start on started

post-start script
  exec su -c /usr/bin/vncserver :1 {{ vnc_user }}
end script
