# roles/vnc/templates/vncserver.j2.conf
description     "start vncserver as {{ default_user_name }}"

start on started

post-start script
  exec su {{ default_user_name }} -c '/usr/bin/vncserver :1'
end script
